# quadlets/agency-inbox.container
[Unit]
Description=Agency inbox monitor (Zoho IMAP reply watcher)
After=agency-outreach.service agency-pod-pod.service
Requires=agency-pod-pod.service

[Container]
Image=localhost/agency-inbox:latest
ContainerName=agency-inbox
Pod=agency-pod.pod
Volume=agency-data:/data:Z
EnvironmentFile=%h/agency/.env

[Service]
Restart=on-failure
TimeoutStartSec=120

[Install]
WantedBy=default.target
