# quadlets/agency-n8n.container
[Unit]
Description=Agency n8n workflow orchestrator
After=agency-pod-pod.service
Requires=agency-pod-pod.service

[Container]
Image=docker.io/n8nio/n8n:latest
ContainerName=agency-n8n
Pod=agency-pod.pod

Volume=agency-n8n-data:/home/node/.n8n:Z
Volume=%h/agency/n8n/workflows:/workflows:Z,ro
EnvironmentFile=%h/agency/.env

Environment=N8N_HOST=0.0.0.0
Environment=N8N_PORT=5678
Environment=N8N_PROTOCOL=http
Environment=DB_TYPE=sqlite
Environment=N8N_BASIC_AUTH_ACTIVE=%E{N8N_BASIC_AUTH_ACTIVE}
Environment=N8N_BASIC_AUTH_USER=%E{N8N_BASIC_AUTH_USER}
Environment=N8N_BASIC_AUTH_PASSWORD=%E{N8N_BASIC_AUTH_PASSWORD}
Environment=N8N_ENCRYPTION_KEY=%E{N8N_ENCRYPTION_KEY}
Environment=WEBHOOK_URL=http://localhost:5678/
Environment=N8N_RUNNERS_ENABLED=true
Environment=GENERIC_TIMEZONE=America/Chicago

[Service]
Restart=on-failure
TimeoutStartSec=120

[Install]
WantedBy=default.target
