# quadlets/agency-marketing.container
[Unit]
Description=Agency marketing agent
After=agency-support.service agency-pod-pod.service
Requires=agency-pod-pod.service

[Container]
Image=localhost/agency-marketing:latest
ContainerName=agency-marketing
Pod=agency-pod.pod
Volume=agency-data:/data:Z
EnvironmentFile=%h/agency/.env

[Service]
Restart=on-failure
TimeoutStartSec=120

[Install]
WantedBy=default.target
