# quadlets/agency-dashboard.container
[Unit]
Description=Agency Streamlit pipeline dashboard
After=agency-pod-pod.service
Requires=agency-pod-pod.service

[Container]
Image=localhost/agency-dashboard:latest
ContainerName=agency-dashboard
Pod=agency-pod.pod

Volume=agency-data:/data:Z
EnvironmentFile=%h/agency/.env


[Service]
Restart=on-failure
TimeoutStartSec=60

[Install]
WantedBy=default.target
